# etc/rc.yenma.in
#!/bin/sh
# $Id$
#
# chkconfig: 2345 98 00
# description: sender-authentication milter
# processname: yenma
# config: @sysconfdir@/yenma.conf
# pidfile: @localstatedir@/@progname@.pid

prefix=@prefix@
exec_prefix=@exec_prefix@
exec=@libexecdir@/yenma
prog=@progname@
config=@sysconfdir@/${prog}.conf
lockfile=/var/lock/subsys/${prog}
NETCAT="/usr/bin/nc"
pid_file=@localstatedir@/${prog}.pid
CONTROL_PORT=`sed -ne "s/^Service.ControlSocket:\s*\(inet6\{0,1\}:\)\{0,1\}\([^@:]*\)\(@.*\)\{0,1\}/\2/p" ${config}`

# soft resource limit of stack size is default size of
# thread specific stack on glibc-2.3.4/glibc-2.4.0
STACK_RLIMIT=512

OPENFILE_RLIMIT=32768

abort () {
    echo "$1" 1>&2
    exit 2
}

control_command () {
    if [ ! -f ${config} -a ! -d ${config} ]; then
        abort "  error: ${config} not found"
    fi

    if [ ! -x ${NETCAT} ]; then
        abort "  error: netcat not found"
    fi

    RESPONSE=`(
        echo "$1"
        sleep 1
        echo "QUIT"
    ) | ${NETCAT} localhost ${CONTROL_PORT}`

    RETVAL=$?

    if [ ${RETVAL} -ne 0 ] ; then
        abort "netcat exited abnormally: exit_code=${RETVAL}"
    fi

    RESPONSE_CODE=`echo "${RESPONSE}" | head -1 | awk '{print $1}'`
    if [ "x${RESPONSE_CODE}" != "x200" ] ; then
        abort "error response from yenma, check log messages. RESPONSE: ${RESPONSE}"
    fi
}

case "$1" in
start)
    if [ ! -x ${exec} ]; then
        abort "  error: ${exec} not found"
    fi

    if [ ! -f ${config} -a ! -d ${config} ]; then
        abort "  error: ${config} not found"
    fi

    ulimit -n ${OPENFILE_RLIMIT}
    ulimit -s ${STACK_RLIMIT}

    ${exec} -c ${config}
    if [ $? -ne 0 ]; then
        abort "  error: yenma failed to start up"
    fi
    touch ${lockfile}
    ;;


stop)
    if [ ! -f ${pid_file} ]; then
        abort "  error: ${pid_file} not found"
    fi

    kill -TERM `cat ${pid_file}`
    if [ $? -ne 0 ]; then
        abort "  error: failed to send SIGTERM to yenma"
    fi
    rm -f ${lockfile}

    echo "  yenma terminating"
    ;;

reload)
    control_command "RELOAD"
    echo "  yenma reloading"
    ;;

graceful-shutdown)
    control_command "GRACEFUL"
    echo "  yenma gracefully shutting down"
    ;;

*)
    echo "Usage: `basename $0` {start|stop|reload|graceful-shutdown}" >&2
    exit 1
    ;;
esac

exit 0
